# dot_config/systemd/user/minikube-autostart.service
[Unit]
Description=Auto-start minikube rootless podman cluster
Documentation=https://minikube.sigs.k8s.io/docs/drivers/podman
# Wait for the network to be up before pulling kicbase/Kubernetes images.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
# The wrapper is in the user's PATH-equivalent location; %h expands to $HOME.
ExecStart=%h/.local/bin/minikube-autostart
# If the start fails (e.g. network not ready yet at boot), retry once after a
# short delay so the cluster is not left down after a transient race.
Restart=on-failure
RestartSec=30

[Install]
# Start when the user manager reaches default.target, which happens at user login
# and — because install-system-config enables lingering — also at boot.
WantedBy=default.target
