# dot_config/systemd/user/readonly_cli-proxy-api.service
[Unit]
Description=CLIProxyAPI managed localhost infrastructure
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=%h/.local/libexec/cli-proxy-api-launch --
Restart=on-failure
RestartSec=2
TimeoutStartSec=15
UMask=0077
Environment=CPA_HOME=%h
Environment=CPA_ACTIVE_BINARY=%h/.local/share/cli-proxy-api/current/cli-proxy-api
Environment=CPA_SOURCE_CONFIG=%h/.config/cli-proxy-api/config.yaml
Environment=CPA_CONFIG=%h/.local/share/cli-proxy-api/runtime/config.yaml
Environment=CPA_AUTH_DIR=%h/.local/share/cli-proxy-api/auth
Environment=CPA_WORK_DIR=%h/.local/share/cli-proxy-api/work
Environment=PATH=/usr/bin:/bin:/usr/sbin:/sbin
Environment=ENV=
Environment=BASH_ENV=
Environment=LD_PRELOAD=
Environment=LD_LIBRARY_PATH=
Environment=LD_AUDIT=
Environment=LD_DEBUG=
Environment=LD_PROFILE=
Environment=GCONV_PATH=
Environment=LOCPATH=
Environment=GLIBC_TUNABLES=
StandardOutput=append:%h/.local/share/cli-proxy-api/work/supervisor.log
StandardError=append:%h/.local/share/cli-proxy-api/work/supervisor.log
# The launcher applies the final environment allowlist and rejects dotenv/auth
# residue before every initial start and automatic restart.

[Install]
WantedBy=default.target
