# system/linux/etc/systemd/system/docker-prune.service
# /etc/systemd/system/docker-prune.service
#
# Oneshot service: prunes unused Docker resources (conservative — no --all).
#
#   docker system prune --force
#     Removes stopped containers, unused networks, *dangling* images
#     (untagged + unreferenced), and build cache. Tagged images that
#     simply have no running container are KEPT. Does NOT touch volumes.
#     https://docs.docker.com/reference/cli/docker/system/prune/
#
#   docker volume prune --force
#     Removes *anonymous* unused volumes. Named volumes are KEPT.
#     https://docs.docker.com/reference/cli/docker/volume/prune/
#
# Fired by docker-prune.timer (weekly). ConditionPathExists gates execution
# on the docker binary being installed, so the unit is safe to ship on hosts
# that haven't run scripts/linux/install-packages.sh yet.

[Unit]
Description=Prune unused Docker resources (dangling images, stopped containers, networks, build cache, anonymous volumes)
Documentation=https://docs.docker.com/reference/cli/docker/system/prune/
Documentation=https://docs.docker.com/reference/cli/docker/volume/prune/
Wants=docker.service
After=docker.service
ConditionPathExists=/usr/bin/docker

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune --force
ExecStart=/usr/bin/docker volume prune --force
