# Library/LaunchAgents/dev.h82.cli-proxy-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>dev.h82.cli-proxy-api</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/sh</string>
		<string>-c</string>
		<string>exec "$HOME/.local/bin/cli-proxy-api" -config "$HOME/.config/cli-proxy-api/config.yaml"</string>
	</array>

	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>2</integer>

	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
