# home/dot_config/systemd/user/mxm4-haptic-notify.service
[Unit]
Description=Desktop-notification to MX Master 4 haptic bridge
# Eavesdrops the session bus for org.freedesktop.Notifications.Notify and
# forwards a waveform to mxm4-hapticd. Needs the session bus, hence the
# graphical-session binding; needs the daemon, hence the ordering + Wants.
After=mxm4-hapticd.service graphical-session.target
Wants=mxm4-hapticd.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/mxm4-haptic-notify
Restart=always
RestartSec=2
# Filtering knobs (uncomment to override defaults):
# Environment=MXM4_NOTIFY_MIN_GAP_MS=800
# Environment=MXM4_NOTIFY_SKIP_REPLACES=1
# Environment=MXM4_NOTIFY_DENY=Spotify,discord

[Install]
WantedBy=graphical-session.target
