# system/linux/etc/udev/rules.d/99-veth-no-ipv6.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="veth*|br-*|podman*", RUN+="/usr/sbin/sysctl -q -w net.ipv6.conf.%k.disable_ipv6=1"
