# dot_config/autostart/private_claude-desktop.desktop
[Desktop Entry]
Type=Application
Name=Claude
Exec="/usr/lib/claude-desktop/claude-desktop" --startup
X-GNOME-Autostart-enabled=true
